# deploy/systemd/whale-web.service.example
[Unit]
Description=Whale CLI WebUI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=whale
Group=whale
WorkingDirectory=/srv/whale-workspace
EnvironmentFile=-/etc/whale-cli.env
Environment=WHALE_HOME=/var/lib/whale-cli
Environment=WHALE_WORKSPACE=/srv/whale-workspace
ExecStart=/opt/whale-cli/.venv/bin/whale-web --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/whale-cli /srv/whale-workspace

[Install]
WantedBy=multi-user.target
